# root/etc/init.d/vlmcsd
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=90

PROG=/usr/bin/vlmcsd

start_service() {
	config_load vlmcsd
	
	local enabled auto_activate internet_access port max_workers lcid activation_interval renewal_interval log_enabled log_verbose
	
	config_get_bool enabled config enabled 0
	[ "$enabled" -eq 1 ] || return 0
	
	config_get auto_activate config auto_activate 1
	config_get internet_access config internet_access 0
	config_get port config port 1688
	config_get max_workers config max_workers 4
	config_get lcid config lcid 1033
	config_get activation_interval config activation_interval "2h"
	config_get renewal_interval config renewal_interval "7d"
	config_get_bool log_enabled config log_enabled 0
	config_get_bool log_verbose config log_verbose 0

	procd_open_instance "vlmcsd"
	procd_set_param command $PROG -D -e
	
	# Listen port
	procd_append_param command -L 0.0.0.0:"$port" -L [::]:"$port"
	
	# Advanced parameters
	[ -n "$max_workers" ] && procd_append_param command -m "$max_workers"
	[ -n "$lcid" ] && procd_append_param command -C "$lcid"
	[ -n "$activation_interval" ] && procd_append_param command -A "$activation_interval"
	[ -n "$renewal_interval" ] && procd_append_param command -R "$renewal_interval"
	
	# Logging
	if [ "$log_enabled" -eq 1 ]; then
		procd_set_param stdout 1
		procd_set_param stderr 1
		[ "$log_verbose" -eq 1 ] && procd_append_param command -v
	fi

	# Fallback to load ini file if user modified it manually
	if [ -s /etc/vlmcsd.ini ]; then
		procd_append_param command -i /etc/vlmcsd.ini
	fi

	procd_set_param respawn
	procd_close_instance
	
	# Auto activate SRV record
	local hostname=$(uci -q get system.@system[0].hostname)
	local srv_rec="_vlmcs._tcp.lan,$hostname.lan,$port,0,100"
	
	local changed=0
	
	# Always clean up any existing _vlmcs records to avoid duplicates/stale ports
	if grep -q "_vlmcs\._tcp\.lan" /etc/config/dhcp; then
		sed -i "/_vlmcs\._tcp\.lan/d" /etc/config/dhcp
		changed=1
	fi
	
	if [ "$auto_activate" = "1" ]; then
		uci add_list dhcp.@dnsmasq[0].srvrecord="$srv_rec"
		changed=1
	fi
	
	if [ "$changed" -eq 1 ]; then
		uci commit dhcp
		/etc/init.d/dnsmasq restart
	fi

	# Firewall rule for Internet access
	if [ "$internet_access" -eq 1 ]; then
		uci -q delete firewall.vlmcsd
		uci set firewall.vlmcsd=rule
		uci set firewall.vlmcsd.name='Allow-vlmcsd'
		uci set firewall.vlmcsd.src='wan'
		uci set firewall.vlmcsd.proto='tcp'
		uci set firewall.vlmcsd.dest_port="$port"
		uci set firewall.vlmcsd.target='ACCEPT'
		uci commit firewall
		/etc/init.d/firewall reload 2>/dev/null
	else
		if uci -q get firewall.vlmcsd >/dev/null; then
			uci -q delete firewall.vlmcsd
			uci commit firewall
			/etc/init.d/firewall reload 2>/dev/null
		fi
	fi
}

service_triggers() {
	procd_add_reload_trigger "vlmcsd"
}

stop_service() {
	if grep -q "_vlmcs\._tcp\.lan" /etc/config/dhcp; then
		sed -i "/_vlmcs\._tcp\.lan/d" /etc/config/dhcp
		uci commit dhcp
		/etc/init.d/dnsmasq restart
	fi
	
	if uci -q get firewall.vlmcsd >/dev/null; then
		uci -q delete firewall.vlmcsd
		uci commit firewall
		/etc/init.d/firewall reload 2>/dev/null
	fi
}
